# rules/systemd.rules
# Monitor systemd service files
-w /etc/systemd/system -p rw -k service
-w /lib/systemd/system -p rw -k service
